# package/DEBIAN/postinst
#!/bin/bash
set -e

CONFIG_FILE="/etc/battery-mqtt-monitor/config"

# ==== INTERACTIVE OR NON-INTERACTIVE ====
if [ -t 0 ]; then
  # Terminal available — ask the user
  echo ""
  echo "=== battery-mqtt-monitor configuration ==="
  echo ""

  read -rp "MQTT Broker IP or hostname [192.168.1.1]: " INPUT_BROKER
  MQTT_BROKER="${INPUT_BROKER:-192.168.1.1}"

  read -rp "MQTT Port [1883]: " INPUT_PORT
  MQTT_PORT="${INPUT_PORT:-1883}"

  read -rp "MQTT Username (leave blank if none): " INPUT_USER
  MQTT_USER="${INPUT_USER:-}"

  if [[ -n "$MQTT_USER" ]]; then
    read -rsp "MQTT Password: " INPUT_PASS
    echo ""
    MQTT_PASS="${INPUT_PASS:-}"
  else
    MQTT_PASS=""
  fi

else
  # No terminal — use defaults and notify user
  echo ""
  echo "⚠ battery-mqtt-monitor: no interactive terminal detected."
  echo "  Default config written to: ${CONFIG_FILE}"
  echo "  Edit it manually, then run:"
  echo "    sudo systemctl restart battery-mqtt-monitor.timer"
  echo ""

  MQTT_BROKER="192.168.1.1"
  MQTT_PORT="1883"
  MQTT_USER=""
  MQTT_PASS=""
fi

# Write config
cat > "$CONFIG_FILE" <<EOF
# battery-mqtt-monitor configuration
# Generated by postinst on $(date +"%Y-%m-%d %H:%M:%S")

MQTT_BROKER="${MQTT_BROKER}"
MQTT_PORT="${MQTT_PORT}"
MQTT_USER="${MQTT_USER}"
MQTT_PASS="${MQTT_PASS}"
EOF

chmod 600 "$CONFIG_FILE"

# Set permissions on script
chmod +x /usr/bin/battery-mqtt-monitor

# Enable and start the timer
systemctl daemon-reload
systemctl enable battery-mqtt-monitor.timer
systemctl start battery-mqtt-monitor.timer

echo ""
echo "✔ battery-mqtt-monitor installed and running."
echo "  Topic:  devices/$(hostname)/battery"
echo "  Timer:  every 60 seconds"
echo "  Config: ${CONFIG_FILE}"
echo ""
